# systemd/watchdog.timer
[Unit]
Description=Run watchdog periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=watchdog.service

[Install]
WantedBy=timers.target
