# systemd/watchdog.service
[Unit]
Description=friday resource watchdog

[Service]
Type=oneshot
ExecStart=/usr/local/bin/watchdog.sh
